# deploy/smap-worker.service
# systemd unit for the scheduler worker. Runs as its own process so web
# restarts don't disturb scheduling.
[Unit]
Description=Social Auto-Poster scheduler worker
After=network.target postgresql.service

[Service]
Type=simple
User=smap
WorkingDirectory=/opt/social-media-auto-posting
EnvironmentFile=/opt/social-media-auto-posting/.env
ExecStart=/opt/social-media-auto-posting/.venv/bin/python -m app.scheduler.worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
